# deployment/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/var/www/deviant/gunicorn.sock

[Install]
WantedBy=sockets.target